# terraform/shared/config/vault_upstart.conf
description "Vault"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

console log

script
  if [ -f "/etc/service/vault" ]; then
    . /etc/service/vault
  fi

  # Make sure to use all our CPUs, because Vault can block a scheduler thread
  export GOMAXPROCS=`nproc`

  exec /usr/local/bin/vault server \
    -config="/etc/vault.d/vault.hcl" \
    \$${VAULT_FLAGS} \
    >>/var/log/vault.log 2>&1
end script
